# docker/php-7.3-lightweight/rootfs/etc/cron.d/logrotate
# Run logrotate daily at 2 AM
0 2 * * * root /usr/sbin/logrotate /etc/logrotate.conf 
